# node/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:sensu-agent]
command=/bin/sensu-agent start
autorestart=true
